# root/Users/ckim/Library/LaunchAgents/ckim.aws-list-instances-mail.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ckim.aws-list-instances-mail</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/aws-list-instances-mail</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/tmp/aws-list-instances-mail-stdout</string>
    <key>StandardErrorPath</key>
    <string>/tmp/aws-list-instances-mail-stderr</string>
    <key>AbandonProcessGroup</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>Nice</key>
    <integer>1</integer>
</dict>
</plist>
